# systemd/retention-manager.timer
[Unit]
Description=Run retention manager every hour

[Timer]
OnCalendar=*-*-* *:30:00
Persistent=true

[Install]
WantedBy=timers.target
